# backup-lvm.cron
1 48 * * * root /root/backup-lvm/lvm-backup.sh
